# packages/ddns-afraid-org/files/etc/systemd/system/ddns-afraid-org.timer
[Unit]
Description=Update DDNS (timer)

[Timer]
OnCalendar=*:0/20:0
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
