# packages/setup-montana/debian/postinst
#!/bin/sh

set -e
case "$1" in
    configure)
        find /etc/ddns-afraid-org.d -type f -exec chmod 640 '{}' \;
        find /etc/ddns-afraid-org.d -type d -exec chmod 750 '{}' \;

        find /root/pki-montana -type d -exec chmod 700 '{}' \;
        find /root/pki-montana -type f -exec chmod 600 '{}' \;
        find /root/pki-montana -type f -name *.sh -exec chmod 700 '{}' \;

        /root/pki-montana/1.first_run.sh

        mkdir -p /etc/nginx/sites-enabled/
        ln -fs /etc/nginx/sites-available/crl  /etc/nginx/sites-enabled/crl
        ln -fs /etc/nginx/sites-available/stat /etc/nginx/sites-enabled/stat

        invoke-rc.d nginx status >/dev/null 2>&1 && invoke-rc.d nginx reload
    ;;
    abort-upgrade|abort-remove|abort-deconfigure)
    ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac
exit 0
