# packages/setup-nevada/etc/systemd/system/montana-crl.service
[Unit]
Description=Download montana-pki CRL (service)
After=network.target network-online.target
Wants=network-online.target

[Service]
ExecStart=if crl=$(curl -s -m 10 --fail-with-body 'http://montana.ajalo.com/crl'); then echo "$crl" > /etc/openvpn/pki-montana/crl.pem; echo "$crl" > /etc/swanctl/x509crl/pki-montana.pem; fi

[Install]
WantedBy=multi-user.target
